# root/etc/services.d/autoscan/run
#!/usr/bin/with-contenv bash

if [ ! -f "$PLEX_AUTOSCAN_CONFIG" ]; then
    s6-svc -d .
    exit 0
fi

until curl -fsS "http://localhost:32400/identity" >/dev/null 2>&1; do sleep 3s; done
sleep 3s

exec /usr/local/bin/autoscan server
